# crontab
*/2 * * * * root date >> /home/registroHora.log
